# systemd/sheyybot-doctor.timer
[Unit]
Description=SheyyBot Doctor - Periodic Health Check
Requires=sheyybot-doctor.service

[Timer]
OnBootSec=60
OnUnitActiveSec=120
AccuracySec=10

[Install]
WantedBy=timers.target
